# systemd/startup-precondition.target
[Unit]
Description=Startup Precondition
Wants=fireplace.service
